# terraform/prd/ansible/roles/dns-server/files/repo-mirror.service
[Unit]
Description=Sync DNF repository mirror

[Service]
Type=oneshot
ExecStart=/usr/local/bin/repo-sync.sh
